# backup/restic/restic-maintenance.timer
[Unit]
Description=Daily restic maintenance (forget + prune)

[Timer]
# Once a day, well after the hourly backups. Keeps prune's expensive API
# usage to a single run per day.
OnCalendar=*-*-* 04:30:00
Persistent=true
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
